# Bash/Linux/wallpaper-switcher.timer
[Unit]
Description=Run wallpaper switcher hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
